# s3/s3backup.timer
[Unit]
Description=Storage s3 test backup timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target