# s3/s3backup.service
[Unit]
Description=Test backup to s3 storage

[Service]
Type=oneshot
ExecStart=/root/backup_s3/backup.sh

[Install]
WantedBy=multi-user.target